# systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (OpenClaw local isolated runtime)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__ROOT__
Environment=HOME=/root
EnvironmentFile=-__ROOT__/.openclaw-home/openclaw-gateway.env
ExecStart=__ROOT__/bin/openclaw-local gateway run --verbose
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=90
SuccessExitStatus=0 143
KillMode=control-group

[Install]
WantedBy=multi-user.target
